# docker/app-backend/laravel-workers.conf
[program:laravel-worker]
directory=/var/www
process_name=%(program_name)s_%(process_num)02d
command=php artisan queue:work --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
numprocs=8
redirect_stderr=true
stdout_logfile=/var/www/storage/logs/worker.log
stopwaitsecs=3600

#[program:app-order-processor]
#directory=/var/www
#process_name=%(program_name)s_%(process_num)02d
#command=php artisan app:order-processor
#autostart=true
#autorestart=true
#redirect_stderr=true
#killasgroup=true
#stdout_logfile=/var/www/storage/logs/order-processor.log
#stopwaitsecs=3600
#user=root
#numprocs=8

[program:php-fpm]
directory=/var/www
process_name=%(program_name)s_%(process_num)02d
command=php-fpm
autostart=true
autorestart=true
redirect_stderr=true
killasgroup=true
stdout_logfile=/var/www/storage/logs/socket-server.log
stopwaitsecs=3600
user=root
numprocs=1



